# unit-files/borg-backup.service
[Unit]
Description=Run the borg backup script.

[Service]
Type=oneshot
User=badger
ExecStart=/home/badger/bin/borg-backup.sh
