# unit-files/borg-backup.timer
[Unit]
Description=Borg backup to NAS. Run's daily.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
