# unit-files/mnt-datarrhea.automount
[Unit]
Description=Automount for datarrhea
ConditionPathExists=/mnt/datarrhea

[Automount]
Where=/mnt/datarrhea
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
